# image/common/systemd/labfoundry-console.service
[Unit]
Description=LabFoundry local appliance console on tty1
Documentation=file:/opt/labfoundry/docs/appliance-console.md
After=systemd-user-sessions.service systemd-networkd.service
Conflicts=getty@tty1.service

[Service]
Type=simple
EnvironmentFile=/etc/labfoundry/labfoundry.env
ExecStart=/opt/labfoundry/.venv/bin/labfoundry-console
Restart=always
RestartSec=2
StandardInput=tty
StandardOutput=tty
StandardError=journal
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
UtmpIdentifier=tty1
UtmpMode=user

[Install]
WantedBy=multi-user.target
